# docker/php-cron/crontab
*/5 * * * * php /var/www/yii apple/cron/recalculate-state